# masscrm_back/docker/php-fpm/crontab
*/1 * * * * /bin/bash /var/www/storage/bash/startQueue.sh > /var/www/work.txt
